# systemd/worktimer.service
[Unit]
Description=Work Timer

[Service]
ExecStart=/usr/bin/worktimer
Restart=always

[Install]
WantedBy=default.target
